# systemd/pidash.service
[Unit]
Description=Pi Digital Dashboard Gateway
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/tacho
EnvironmentFile=/home/pi/tacho/.env
ExecStart=/home/pi/tacho/.venv/bin/python /home/pi/tacho/run.py
Restart=always
RestartSec=2
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
