# systemd/pidash-kiosk.service
[Unit]
Description=Pi Dashboard Kiosk Screen
After=pidash.service systemd-user-sessions.service
Requires=pidash.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/tacho
Environment=HOME=/home/pi
Environment=DASH_URL=http://127.0.0.1:8080
ExecStart=/home/pi/tacho/scripts/start-kiosk.sh
Restart=always
RestartSec=2
StandardInput=tty
TTYPath=/dev/tty1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target